# ops/systemd/surfai-backup.service
[Unit]
Description=SURFAI postgres daily backup
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
User=root
Group=root
# Ensure the backup directory exists BEFORE systemd sets up the mount
# namespace for ReadWritePaths. The `+` prefix runs this command outside
# the sandbox with full privileges — required because the sandbox is
# what fails when the path doesn't exist yet.
ExecStartPre=+/bin/mkdir -p /opt/surfai/backups
ExecStart=/opt/surfai/ops/backup.sh
TimeoutStartSec=900
StandardOutput=journal
StandardError=journal
SyslogIdentifier=surfai-backup

# Hardening: only the backup dir is writable.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/surfai/backups /tmp
PrivateTmp=true
